# systemd/fip-screen.service
[Unit]
Description=FIP now-playing screen (pygame)
After=lightdm.service

[Service]
User=@USER@
Group=@USER@
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/@UID@
Environment=SDL_VIDEODRIVER=wayland
ExecStart=/opt/fipradio/venv/bin/python3 /opt/fipradio/pygame_ui.py
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
